# ansible/roles/gpaste/templates/gpaste.service.j2
[Unit]
Description=Clipboard manager client

[Service]
ExecStart=/usr/bin/gpaste-client start

[Install]
WantedBy=multi-user.target
